# assets/all/ld.so.preload
/support/libdisableselinux.so
/support/libandroidshmem.so
