# assets/all/ld.so.preload
/support/libdisableselinux.so
/support/libandroid-shmem.so
